# ops/launchd/poly_strategy_realtime_discovery_refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>poly_strategy_realtime_discovery_refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/ww/Project/poly_strategy &amp;&amp; PYTHONUNBUFFERED=1 LIMIT=100 PAGES=5 BATCH_SIZE=5 CONTEXT_MARKET_LIMIT=0 LLM_WORKERS=4 LLM_ERROR_RETRIES=2 LLM_ERROR_RETRY_BATCH_SIZE=1 LLM_FALLBACK_MODEL=gpt-5.5 LLM_FALLBACK_ERROR_RETRIES=2 LLM_FALLBACK_ERROR_RETRY_BATCH_SIZE=1 LLM_TIMEOUT=180 LLM_REASONING_EFFORT=high INCLUDE_TOP_MARKETS=150 INCLUDE_TOP_NEG_RISK_GROUPS=25 MAX_WATCHLIST_MARKETS=250 scripts/refresh_discovery_watchlist.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>WorkingDirectory</key>
  <string>/Users/ww/Project/poly_strategy</string>
  <key>StandardOutPath</key>
  <string>/Users/ww/Project/poly_strategy/data/realtime-discovery-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ww/Project/poly_strategy/data/realtime-discovery-refresh.log</string>
</dict>
</plist>
