# ops/launchd/poly_strategy_realtime_monitor_24h.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>poly_strategy_realtime_monitor_24h</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/ww/Project/poly_strategy &amp;&amp; PYTHONUNBUFFERED=1 REPORT_OUT=data/realtime-monitor-24h-v1.jsonl SNAPSHOTS_OUT=data/realtime-monitor-24h-v1-snapshots.ndjson SNAPSHOT_INTERVAL=2 STALE_TIMEOUT=30 RECONNECT_DELAY=2 MIN_NET_EDGE=0.002 MAX_CAPITAL_PER_TRADE=10 BANKROLL=50 MIN_PAPER_ROI=0.01 MIN_RUN_OBSERVATIONS=2 MIN_RUN_SECONDS=3 INCLUDE_TOP_MARKETS=150 INCLUDE_TOP_NEG_RISK_GROUPS=25 MAX_WATCHLIST_MARKETS=250 scripts/run_realtime_monitor.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/ww/Project/poly_strategy</string>
  <key>StandardOutPath</key>
  <string>/Users/ww/Project/poly_strategy/data/realtime-monitor-24h-v1.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ww/Project/poly_strategy/data/realtime-monitor-24h-v1.log</string>
</dict>
</plist>
